# crontab
MAILTO=root

42 02 * * 2 root /usr/sbin/letsencrypt.sh -c > /dev/null
