# crontab
MAILTO=root

42 02 1 * * /bin/sh -c "/usr/sbin/letsencrypt.sh -c `[ -x /etc/webapps/letsencrypt.sh/hook.sh ] && echo '--hook /etc/webapps/letsencrypt.sh/hook.sh'` > /dev/null"
