# crontab
MAILTO=root

42 02 1 * * /usr/sbin/letsencrypt.sh -c > /dev/null
