# letsencrypt.sh.crontab
# every month
0 0 1 * * /usr/sbin/letsencrypt.sh -c
